# supervisor.d/kill_inactive_users.conf
[program:kill_users]
command=/usr/bin/kill_inactive_users.py
user=root
autorestart=true
priority=100
